# systemd-recipe-app.service
[Unit]
Description=RasPi Recipe App
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/recipe-app
ExecStart=/usr/local/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
EnvironmentFile=/home/pi/recipe-app/.env

[Install]
WantedBy=multi-user.target
